# packaging/debian/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    systemd-sysusers /usr/lib/sysusers.d/vstimd.conf

    # Warn if a display manager is enabled — it holds DRM master and will
    # prevent vstimd from acquiring the display via VK_KHR_display.
    for dm in gdm gdm3 lightdm sddm xdm; do
        if systemctl is-enabled --quiet "$dm" 2>/dev/null; then
            echo "vstimd: WARNING: $dm is enabled and will conflict with vstimd."
            echo "  Disable it with: sudo systemctl disable --now $dm"
        fi
    done
fi

#DEBHELPER#
