# gpiochip-daqd/packaging/systemd/gpiochip-daqd.service
[Unit]
Description=Linux gpiochip ↔ VTL trigger-line bridge
Documentation=https://github.com/braemons/vstimd
# Start after vstimd so the VTL shared-memory segment exists.
# BindsTo ensures we stop if vstimd exits (segment would be unlinked).
After=vstimd.service
BindsTo=vstimd.service

[Service]
Type=notify
ExecStart=/usr/bin/gpiochip-daqd
Restart=on-failure
RestartSec=2s

# GPIO character device requires root or membership in the gpio group.
User=root
SupplementaryGroups=gpio

# Logging — forward to rsyslog via journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpiochip-daqd

# Real-time scheduling: set SCHED_FIFO before exec so all threads inherit it.
# The Rust code then raises per-thread priorities (output=60, input=50).
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
CPUSchedulingResetOnFork=no

# CAP_SYS_NICE is required so threads can adjust their own SCHED_FIFO priority
# after spawning (output thread raises to 60).
CapabilityBoundingSet=CAP_SYS_NICE
AmbientCapabilities=CAP_SYS_NICE

[Install]
WantedBy=multi-user.target
